# cookbooks/cluster-mesos/templates/default/consul_client.conf.erb
description "Consul client process"
start on (local-filesystems and net-device-up IFACE=eth0)
stop on runlevel [!12345]
respawn
exec sudo consul agent -config-dir /etc/consul.d/client
